# app/camera/config/camera-ota-installer.path
[Unit]
Description=Watch for camera OTA install triggers

[Path]
# Fires when the trigger file is created/modified. The unprivileged
# camera-streamer process writes this file after staging a bundle,
# which activates camera-ota-installer.service (runs as root).
PathModified=/var/lib/camera-ota/trigger
Unit=camera-ota-installer.service

[Install]
WantedBy=multi-user.target
