# app/camera/config/camera-ota-reboot.path
[Unit]
Description=Watch for camera OTA reboot triggers

[Path]
# Fires when the trigger file is created/modified. The unprivileged
# camera-streamer process writes this file from
# `/api/ota/reboot` after a successful install (status=installed).
# This activates camera-ota-reboot.service (which runs as root and
# calls `systemctl reboot`).
PathModified=/var/lib/camera-ota/reboot-trigger
Unit=camera-ota-reboot.service

[Install]
WantedBy=multi-user.target
