# app/camera/config/ensure-camera-overlay.service
# REQ: SWR-053, SWR-050; RISK: RISK-007, RISK-022; SEC: SC-019; TEST: TC-012, TC-047
[Unit]
Description=Ensure camera sensor overlay in config.txt
DefaultDependencies=no
After=local-fs.target
Before=camera-streamer.service

[Service]
Type=oneshot
ExecStart=/opt/camera/scripts/ensure-camera-overlay.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
