# meta-home-monitor/recipes-multimedia/mediamtx/files/mediamtx.service
# REQ: SWR-006, SWR-031, SWR-050; RISK: RISK-001, RISK-017; SEC: SC-016, SC-019; TEST: TC-001, TC-028, TC-044
[Unit]
Description=MediaMTX RTSP Server
After=network-online.target
Wants=network-online.target
Before=monitor.service

[Service]
Type=simple
ExecStart=/usr/bin/mediamtx /etc/mediamtx/mediamtx.yml
Restart=always
RestartSec=3

# Security hardening
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
